# packaging/systemd/cubby.service.tmpl
# Reference systemd user service. `cubby install` generates this automatically;
# this copy is here for manual setup. Replace {{CUBBY_BIN}}, drop it in
# ~/.config/systemd/user/cubby.service, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now cubby.service

[Unit]
Description=Cubby - sort the Downloads folder
After=default.target

[Service]
Type=simple
ExecStart={{CUBBY_BIN}} watch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
